# .config/autostart/nvidia-startup.desktop
 [Desktop Entry]
Type=Application
Name=nvidia-script
Exec=/home/sachnr/.local/bin/nvidia-startup
OnlyShowIn=GNOME;
